# supervisor.ini
[program:shop_back]
command=uwsgi --ini /home/yang/myshop-back/uwsgi.ini ;程序启动
priority=1 ; 数字越高，优先级越高
numprocs=1 ; 启动几个进程
autostart=true ; 随着supervisord的启动而启动
autorestart=true ; 自动重启
startretries=10 ; 启动失败时的最多重试次数
redirect_stderr=true ; 重定向stderr到stdout
stdout_logfile = /home/yang/myshop-back/shop_front.log

[program:shop_api]
command=gunicorn -c /home/yang/myshop-api/gunicorn.py myshop.wsgi:application ;程序启动
directory=/home/yang/myshop-api
priority=1 ; 数字越高，优先级越高
numprocs=1 ; 启动几个进程
autostart=true ; 随着supervisord的启动而启动
autorestart=true ; 自动重启
startretries=10 ; 启动失败时的最多重试次数
redirect_stderr=true ; 重定向stderr到stdout
stdout_logfile = /home/yang/myshop-api/shop_api.log